# rpi/service/log/run
#!/bin/sh
exec setuidgid pi multilog /var/log/sync2serial
