# rpi/service/log/run
#!/bin/sh
exec /usr/bin/envuidgid pi /usr/bin/multilog /var/log/sync2serial
