# rpi/service/log/run
#!/bin/sh
exec multilog /var/log/sync2serial
